# deploy/macos/com.arcanada.disk-arcana.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  DISK-0006 R11 — launchd LaunchDaemon for Disk Arcana client.

  Install path: /Library/LaunchDaemons/com.arcanada.disk-arcana.plist
  Owner: root:wheel, mode 0644 (launchd refuses world-writable plists).
  Load: sudo launchctl bootstrap system /Library/LaunchDaemons/com.arcanada.disk-arcana.plist
  Unload: sudo launchctl bootout system /Library/LaunchDaemons/com.arcanada.disk-arcana.plist
  Status: sudo launchctl print system/com.arcanada.disk-arcana
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.arcanada.disk-arcana</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/disk</string>
        <string>daemon</string>
        <string>start</string>
        <string>--foreground</string>
        <string>--config</string>
        <string>/etc/disk-arcana/disk.toml</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>StandardOutPath</key>
    <string>/var/log/disk-arcana/disk-daemon.out.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/disk-arcana/disk-daemon.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
        <key>HOME</key>
        <string>/var/lib/disk-arcana</string>
    </dict>

    <!-- DISK-0057 (P5-S4b): Mac is single-user; the daemon syncs into the operator's
         git tree. Running as root creates root-owned files causing excessive git churn
         (every synced file appears modified due to ownership change). Run as the operator.
         Replace __OPERATOR_USER__ with the actual login name at install time:
           sudo sed "s/__OPERATOR_USER__/$(id -un)/g" \
             /Library/LaunchDaemons/com.arcanada.disk-arcana.plist > /tmp/plist.tmp && \
           sudo cp /tmp/plist.tmp /Library/LaunchDaemons/com.arcanada.disk-arcana.plist
         Or use the install script: scripts/install-macos.sh (substitutes automatically).
         See runbook: documentation/runbooks/disk-arcana-prod-sync.md § P5-R-S4 -->
    <key>UserName</key>
    <string>__OPERATOR_USER__</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>WorkingDirectory</key>
    <string>/var/lib/disk-arcana</string>
</dict>
</plist>
